# overlay/usr/local/etc/rc.d/adguard
#!/bin/sh
#
# PROVIDE: adguard
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# adguard_enable (bool):        Set to NO by default.
#                       Set it to YES to enable it.
# adguard_user:         The user account adguard daemon runs as what
#                       you want it to be. It uses 'root' user by
#                       default.
# adguard_group:        The group account adguard daemon runs as what
#                       you want it to be. It uses 'wheel' group by
#                       default.
# adguard_pidfile:      Set the location of the adguard pidfile
# adguard_logfile:      Set the location of the adguard logfile
# adguard_workdir:      Set the location of the adguard workdir

. /etc/rc.subr

name=adguard
rcvar=adguard_enable
load_rc_config ${name}

: ${adguard_enable:=NO}
: ${adguard_user:=root}
: ${adguard_group:=wheel}
: ${adguard_workdir:="/usr/local/adguard"}
: ${adguard_pidfile:="/var/run/adguard/adguard.pid"}
: ${adguard_logfile:="/var/log/adguard/adguard.log"}

pidfile=${adguard_pidfile}
logfile=${adguard_logfile}
workdir=${adguard_workdir}

command=/usr/sbin/daemon
procname="/root/AdGuardHome/AdGuardHome"
command_args="-f ${procname} --host 0.0.0.0 --port 3000 --pidfile ${pidfile} --logfile ${logfile} --work-dir ${workdir}"

start_precmd="${name}_prestart"

adguard_prestart()
{
	PATH=${PATH}:/usr/local/bin:/usr/local/sbin
	export LC_CTYPE="en_US.UTF-8"
	for dir in ${workdir} ${logfile%/*} ${pidfile%/*}; do
		if [ ! -d "${dir}" ]; then
				install -d -o ${adguard_user} -g ${adguard_group} ${dir}
		fi
	done
}

run_rc_command "$1"
